# virtualbox-qt5/vmsvga.service
[Unit]
Description=VMSVGA Guest Service

[Service]
User=root
ExecStartPre=-/usr/bin/VBoxClient --vmsvga

[Install]
WantedBy=multi-user.target
