# virtualbox-qt5/vmsvga.service
[Unit]
Description=VMSVGA Guest Service
BindsTo=sddm.service

[Service]
ExecStart=/usr/bin/VBoxDRMClient

[Install]
WantedBy=multi-user.target
