# virtualbox-qt5/vmsvga.service
[Unit]
Description=VMSVGA Guest Service
BindsTo=sddm.service

[Service]
ExecStartPre=-/usr/bin/VBoxClient --vmsvga
ExecStart=/usr/bin/whoami

[Install]
WantedBy=multi-user.target
